# launchd/com.daily-driver.checkin.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.daily-driver.checkin</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>PLACEHOLDER_REPO/scripts/open-session.sh</string>
    <string>check-in</string>
  </array>
  <key>StartCalendarInterval</key>
  <array>
PLACEHOLDER_CHECKIN_INTERVALS
  </array>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>PLACEHOLDER_REPO/logs/checkin.log</string>
  <key>StandardErrorPath</key>
  <string>PLACEHOLDER_REPO/logs/checkin.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>PLACEHOLDER_BREW/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>PLACEHOLDER_HOME</string>
  </dict>
</dict>
</plist>
